# conf/supervisord.conf
[supervisord]
logfile=/root/logs/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/root/logs/supervisord/               ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
 
[program:uwsgi]
environment=PATH=/root/DAPOS_corpus_browser/env/bin
directory=/root/DAPOS_corpus_browser/conf
command=uwsgi --ini uwsgi.ini
autorestart=true

[program:celeryd]
environment=PATH=/root/DAPOS_corpus_browser/env/bin
directory=/root/DAPOS_corpus_browser/corpus_browser/
command=python manage.py celeryd
autorestart=true

[program:celery_beat]
environment=PATH=/root/DAPOS_corpus_browser/env/bin
directory=/root/DAPOS_corpus_browser/corpus_browser/
command=python manage.py celeryd -B
autorestart=true

[program:cat]
environment=PATH=/root/DAPOS_corpus_browser/env/bin
directory=/root/DAPOS_corpus_browser/
command=memcached -vv
autorestart=true
